# config/profile.d/alias.sh
alias docker='podman' 2>/dev/null
alias sudo='doas' 2>/dev/null
alias bcroot='doas su - root' 2>/dev/null
