# config/cron.daily/apparmor-parser
#!/bin/sh
for profile in /etc/apparmor.d/*; do
    [ -f "$profile" ] && /usr/bin/apparmor_parser --replace "$profile"
done
