# service/agetty-tty8-autologin-admin/run
#!/bin/sh
exec agetty --autologin admin --noclear tty8 38400 linux
